# ansible/roles/lego/templates/lego@.timer.j2
[Unit]
Description=Run lego

[Timer]
OnCalendar={{ lego_systemd_calendar }}
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=timers.target