# ansible/roles/runrestic/templates/runrestic@.timer.j2
[Unit]
Description=Run runrestic

[Timer]
OnCalendar={{ runrestic_systemd_calendar }}
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target