# ansible/roles/dnsdist/files/dnsdist.service
[Service]
ExecStart=
ExecStart=/usr/bin/dnsdist \
  --supervised \
  --disable-syslog \
  -C /etc/dnsdist/dnsdist.yml
