# ansible/roles/powerdns_recursor/files/pdns-recursor-rpz-sync.timer
[Unit]
Description=pdns recursor rpz sync cron

[Timer]
OnCalendar=daily
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=timers.target